# deploy/roles/cse191/files/cse191web.conf
description "CSE191 web server instance"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec uwsgi /home/cse191/cse191/web/uwsgi.ini
